# deploy/debian13/systemd/staysuite-realtime.service
# =============================================================================
# StaySuite HospitalityOS - Realtime WebSocket Service
# =============================================================================
# Handles real-time notifications, booking updates, and live data sync
# via Socket.IO on port 3003
#
# Install:
#   sudo cp staysuite-realtime.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now staysuite-realtime
#
# Manage:
#   sudo systemctl restart staysuite-realtime
#   journalctl -u staysuite-realtime -f
# =============================================================================

[Unit]
Description=StaySuite HospitalityOS - Realtime WebSocket Service (port 3003)
Documentation=https://github.com/staysuite/hospitality-os
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=staysuite
Group=staysuite
WorkingDirectory=/opt/staysuite/mini-services/realtime-service
EnvironmentFile=/opt/staysuite/.env
ExecStart=/usr/local/bin/bun --hot index.ts
PIDFile=/opt/staysuite/mini-services/realtime-service/.realtime.pid
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

# --- Security Hardening ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/staysuite/mini-services/realtime-service /tmp
PrivateTmp=true

# --- Resource Limits ---
LimitNOFILE=65535
LimitNPROC=1024
MemoryHigh=512M
MemoryMax=1G

# --- Logging ---
StandardOutput=journal
StandardError=journal
SyslogIdentifier=staysuite-realtime

[Install]
WantedBy=multi-user.target
